# Jarvis/jarvis.service
[Unit]
Description=JARVIS AI Assistant Server
After=network.target

[Service]
Type=simple
User=Grindlewalt
WorkingDirectory=/mnt/storage/NAS/Jarvis/Jarvis
ExecStart=/mnt/storage/NAS/Jarvis/.venv/bin/python3 /mnt/storage/NAS/Jarvis/Jarvis/run.py
Restart=always
RestartSec=10
StandardOutput=append:/mnt/storage/NAS/Jarvis/Jarvis/logs/jarvis.log
StandardError=append:/mnt/storage/NAS/Jarvis/Jarvis/logs/jarvis.log
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
