# jarvis/jarvis.service
[Unit]
# Jarvis personal AI assistant — see /mnt/storage/NAS/Jarvis/jarvis/docs/BUILD_SPEC.md
Description=Jarvis personal AI assistant
After=network-online.target ollama.service
Wants=network-online.target

[Service]
Type=simple
User=Grindlewalt
Group=Grindlewalt

# Per BUILD_SPEC §1: WorkingDirectory=/mnt/storage/NAS/Jarvis/jarvis. The package
# itself lives at this path (this dir IS the `jarvis` Python package). To make
# `python -m jarvis.run` resolve, we put the parent on PYTHONPATH so import
# discovery finds `jarvis/` as a top-level package.
WorkingDirectory=/mnt/storage/NAS/Jarvis/jarvis
Environment="PYTHONPATH=/mnt/storage/NAS/Jarvis"
Environment="PYTHONUNBUFFERED=1"

# AGENT_CENTRAL_MIRROR_OWNER=jarvis flips Jarvis to the mirror writer.
# CMD-Claude shipped the env-flag check prior to P9; with this set CMD's
# own (now-disabled) curator backs off, leaving Jarvis the sole writer of
# ~/.agent_bin/central_context.md (BUILD_SPEC §15, §21.10).
Environment="AGENT_CENTRAL_MIRROR_OWNER=jarvis"

ExecStart=/mnt/storage/NAS/Jarvis/.venv/bin/python -m jarvis.run

# P5 mounts the FastAPI server in the main thread; the daemon is now
# long-lived. Restart on failure (not on clean exit, which is a SIGTERM).
Restart=on-failure
RestartSec=5

# Resource hygiene
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
